# debian/cron.d
# mins hour dofm mnth dofw   user          command
# ---- ---- ---- ---- ----   -----------   ----------------------------------------------------------------------------
     *    *    *    *    *   root          dropfile
# ---- ---- ---- ---- ----   -----------   ----------------------------------------------------------------------------
